# root/etc/s6-overlay/s6-rc.d/pty-server/finish
#!/command/with-contenv bash
# s6 finish script for pty-server
exit 0
